# scripts/meshsig.service
[Unit]
Description=MeshSig — Cryptographic security layer for AI agents
After=network.target

[Service]
Type=simple
# Standalone: dashboard + API on port 4888
# Proxy mode: add --gateway http://localhost:YOUR_GATEWAY_PORT
# and set MESH_PORT to the port your agents call
ExecStart=/usr/bin/node /opt/meshsig/dist/main.js start --no-terminal
Restart=always
RestartSec=5
User=meshsig
Group=meshsig
Environment=MESH_PORT=4888
Environment=NODE_ENV=production
WorkingDirectory=/opt/meshsig
# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meshsig

[Install]
WantedBy=multi-user.target
